# worker-service/crontab
PATH=/usr/local/bin:/usr/bin:/bin
41 22 * * * source /etc/backup.env && cd /django && python3 manage.py backup_db
